# templates/token-renewal.service.erb
[Unit]
Description=vault token renewal service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
EnvironmentFile=<%= scope['vault_client::config_path'] %>
ExecStart=<%= scope['vault_client::bin_dir'] %>/<%= scope['vault_client::app_name'] %> renew-token --init-role=${VAULT_INIT_ROLE}

[Install]
WantedBy=multi-user.target
